# deploy/litestream.example.service
[Unit]
Description=Litestream SQLite Replication
After=network.target

[Service]
Type=simple
User=<user>
ExecStart=/usr/local/bin/litestream replicate -config /home/<user>/ordinary-note/deploy/litestream.yml
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
