# deploy/litestream.example.service
[Unit]
Description=Litestream SQLite Replication
After=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/litestream replicate -config /home/<user>/ordinary-note/deploy/litestream.yml
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
